# deploy/ktta.service
# KTTA Platform - systemd service unit
# /etc/systemd/system/ktta.service にコピー
# sudo systemctl daemon-reload && sudo systemctl enable --now ktta

[Unit]
Description=KTTA Platform (Kushiro Table Tennis Association)
Documentation=https://github.com/AKIHIKO-S/kushiro
After=network.target

[Service]
Type=simple
User=ktta
WorkingDirectory=/opt/ktta
ExecStart=/usr/bin/node /opt/ktta/server.js
Restart=always
RestartSec=3
Environment=NODE_ENV=production
Environment=PORT=3000
# 管理キー + DB パス等は /etc/ktta.env から読み込む (chmod 600)
EnvironmentFile=-/etc/ktta.env
StandardOutput=journal
StandardError=journal
SyslogIdentifier=ktta

# セキュリティ強化
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
# /opt/ktta = アプリ本体, /var/data = SQLite DB と印鑑画像
ReadWritePaths=/opt/ktta /var/data
ProtectHome=true

# リソース制限 (任意・Always Free A1 でも余裕)
MemoryHigh=512M
MemoryMax=1G

[Install]
WantedBy=multi-user.target
